# configs/linux.timer
[Unit]
Description=Timer for Scheduled Shutdown/Restart

[Timer]
OnCalendar=*-*-* 22:00:00  # Updated dynamically

[Install]
WantedBy=timers.target
